# deploy/roles/project/templates/celery.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!
# {{ansible_managed}}

[program:worker-{{base_project_name}}]
autorestart=true
autostart=true
command={{wsgi_virtualenv}}/bin/celery worker -A {{base_project_name}}.main
directory={{base_source_directory}}
redirect_stderr=true
stdout_logfile={{base_logs_directory}}/worker.log
user={{base_deploy_user}}

[program:beat-{{ base_project_name }}]
autorestart=true
autostart=true
command={{wsgi_virtualenv}}/bin/celery beat -A {{base_project_name}}.main
directory={{base_source_directory}}
redirect_stderr=true
stdout_logfile={{base_logs_directory}}/beat.log
user={{base_deploy_user}}
